# example/supervisor.ini

[program:server1]
directory = /tmp/
command=python /tmp/server_1.py
autostart=false
autorestart=true

[program:server2]
directory = /tmp/
command=python /tmp/server_2.py
autostart=false
autorestart=true
